# infra/wslcb-task@disk-hygiene.service.d/override.conf
[Service]
# disk-hygiene is genuinely idle-priority background work; scoped to just
# this instance so scrape/backfill-addresses keep default scheduling priority.
Nice=10
IOSchedulingClass=idle
