# wslcb-task@.service
[Unit]
Description=WSLCB task: %i
After=network.target

[Service]
Type=oneshot
User=exedev
WorkingDirectory=/home/exedev/wslcb-licensing-tracker
ExecStart=/home/exedev/wslcb-licensing-tracker/venv/bin/python scraper.py %i
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/exedev/wslcb-licensing-tracker/env
# Default TimeoutStartSec for oneshot is infinity on most distros;
# if this changes, set it explicitly here (refresh can take ~20 min).
