# wslcb-web.service
[Unit]
Description=WSLCB Licensing Tracker Web App
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/wslcb-tracker
ExecStart=/home/exedev/wslcb-tracker/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
