# wslcb-scraper.timer
[Unit]
Description=Run WSLCB scraper twice daily (12:30 AM and 6:30 AM Pacific)

[Timer]
# The daily scrape uses the wslcb-task@ template with "scrape" instance
Unit=wslcb-task@scrape.service
OnCalendar=*-*-* 00:30:00 America/Los_Angeles
OnCalendar=*-*-* 06:30:00 America/Los_Angeles
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
